# iso-builder/configs/system/autologin.service
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin hypr --noclear %I $TERM
